# gauge.service
[Unit]
Description=CHT Gauge - Round SPI Display
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
WorkingDirectory=/root/gauge
# Force the kernel framebuffer driver (required for smooth updates)
Environment=SDL_VIDEODRIVER=dummy
Environment=SDL_FBDEV=/dev/fb0
Environment=PYTHONUNBUFFERED=1

# Best-effort console cleanup so nothing else is fighting the fb
ExecStartPre=-/bin/sh -c '
  chvt 1 2>/dev/null || true
  echo 0 > /sys/class/graphics/fbcon/cursor_blink 2>/dev/null || true
  systemctl stop getty@tty1.service serial-getty@ttyS0.service serial-getty@ttyAMA0.service 2>/dev/null || true
'

ExecStart=/root/gauge/venv/bin/python3 /root/gauge/gauge.py

# Keep it alive
Restart=always
RestartSec=5

# All output (including any errors) goes to the journal.
# The actual display output is done by direct writes to /dev/fb0 in the Python code.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
