# gauge.service
[Unit]
Description=CHT Gauge - Round SPI Display
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
WorkingDirectory=/root/gauge
# Force the kernel framebuffer driver (required for smooth updates)
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_FBDEV=/dev/fb0
Environment=PYTHONUNBUFFERED=1

# Aggressively take over the console (stop any getty that might be on tty1 or serial)
ExecStartPre=-/bin/sh -c '
  chvt 1 2>/dev/null || true
  echo 0 > /sys/class/graphics/fbcon/cursor_blink 2>/dev/null || true
  systemctl stop getty@tty1.service serial-getty@ttyS0.service serial-getty@ttyAMA0.service 2>/dev/null || true
'

ExecStart=/root/gauge/venv/bin/python3 /root/gauge/gauge.py

# Keep it alive if it crashes (it will retry every few seconds)
Restart=always
RestartSec=5

# Log to journal (view with journalctl -u gauge.service).
# This prevents error text from being painted directly onto the framebuffer as "terminal" output.
# When the gauge successfully initializes, it will draw its graphics over the fb.
StandardOutput=journal
StandardError=journal
# We still try to attach to the right tty for the driver
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
